# Filesystem/rootfs_files/etc/profile.d/qt5.sh
#insmod /lib/modules/galcore.ko
/root/qt_example/mainwindow &
